# systemd/kble-serialport.service
[Unit]
Description=kble-serialport (UART <-> WebSocket bridge)
After=network-online.target
Wants=network-online.target
# Allow up to 5 restarts within 60s before giving up; without this a crash-loop
# would hit the default start limit and leave the unit permanently stopped.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
# --addr is the Twingate-reachable interface IP (or 0.0.0.0 + host firewall).
# --port is the TCP listen port; clients pick the serial device via ?port= in the URL.
ExecStart=@KBLE@ --addr @ADDR@ --port 9600
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
